# systemd/ssh-tunnel.service.example
[Unit]
Description=Reverse SSH Tunnel to VPS
After=network.target

[Service]
Type=simple
User=wir
Environment="AUTOSSH_GATETIME=0"
ExecStart=/usr/bin/autossh -M 0 -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -o BatchMode=yes -i /home/wir/.ssh/telegram_agent -R 127.0.0.1:18080:127.0.0.1:8081 USER@VPS_HOST
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
